# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize --fpm-config /usr/local/etc/php-fpm.d/www.conf
autostart=true
autorestart=true
startretries=3
startsecs=1
user=www-data
stdout_logfile=/var/log/php/php-fpm.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/php/php-fpm-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
environment=LANG=en_US.UTF-8

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
startsecs=1
user=www-data
stdout_logfile=/var/log/nginx/nginx.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/nginx/nginx-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
environment=LANG=en_US.UTF-8

[program:queue-worker]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
startretries=3
startsecs=1
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/var/log/php/queue-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/php/queue-worker-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
environment=LANG=en_US.UTF-8

[program:schedule-run]
command=php /var/www/html/artisan schedule:run
autostart=true
autorestart=true
startretries=3
startsecs=1
user=www-data
stdout_logfile=/var/log/php/schedule.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/php/schedule-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
environment=LANG=en_US.UTF-8

[program:horizon]
command=php /var/www/html/artisan horizon
autostart=true
autorestart=true
startretries=3
startsecs=1
user=www-data
stdout_logfile=/var/log/php/horizon.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/php/horizon-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
environment=LANG=en_US.UTF-8

[program:redis]
command=redis-server --requirepass %(ENV_REDIS_PASSWORD)s --appendonly yes --appendfsync everysec
autostart=true
autorestart=true
startretries=3
startsecs=1
user=redis
stdout_logfile=/var/log/redis/redis.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/redis/redis-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
environment=LANG=en_US.UTF-8

[program:backup]
command=/bin/sh -c "while true; do sleep 86400; php /var/www/html/artisan backup:run; done"
autostart=true
autorestart=true
startretries=3
startsecs=1
user=www-data
stdout_logfile=/var/log/php/backup.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/php/backup-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
environment=LANG=en_US.UTF-8

[program:monitor]
command=php /var/www/html/artisan monitor:run
autostart=true
autorestart=true
startretries=3
startsecs=1
user=www-data
stdout_logfile=/var/log/php/monitor.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/php/monitor-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
environment=LANG=en_US.UTF-8

[program:logrotate]
command=/usr/sbin/logrotate /etc/logrotate.conf
autostart=true
autorestart=true
startretries=3
startsecs=1
user=root
stdout_logfile=/var/log/logrotate.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/logrotate-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
environment=LANG=en_US.UTF-8

[eventlistener:processes]
command=/usr/bin/supervisor-event-listener
events=PROCESS_STATE,TICK_5
buffer_size=1024
events=PROCESS_STATE,TICK_5
stdout_logfile=/var/log/supervisor/eventlistener.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/eventlistener-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

[group:app]
programs=php-fpm,nginx,queue-worker,schedule-run,horizon,redis,backup,monitor,logrotate